# deploy/systemd/fininzen-refresh-prices.service
[Unit]
Description=Fininzen asset price refresh
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=fininzen
Group=fininzen
WorkingDirectory=/opt/fininzen
EnvironmentFile=/etc/fininzen.env
ExecStart=/opt/fininzen/venv/bin/python manage.py refresh_asset_prices
Nice=10
IOSchedulingClass=idle
NoNewPrivileges=yes
PrivateTmp=yes
ProtectHome=yes
ProtectSystem=full
ReadWritePaths=/opt/fininzen
