# deploy/systemd/finnet-refresh-prices.service
[Unit]
Description=Finnet asset price refresh
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=finnet
Group=finnet
WorkingDirectory=/opt/finnet
EnvironmentFile=/etc/finnet.env
ExecStart=/opt/finnet/venv/bin/python manage.py refresh_asset_prices
Nice=10
IOSchedulingClass=idle
NoNewPrivileges=yes
PrivateTmp=yes
ProtectHome=yes
ProtectSystem=full
ReadWritePaths=/opt/finnet
